# f5-ast-logic/rootfs/etc/services.d/init-remote-db/run
#!/usr/with-contenv bashio
HOST=$(bashio::config 'db_host'); PORT=$(bashio::config 'db_port')
USER=$(bashio::config 'db_user'); PASS=$(bashio::config 'db_password')

: ${HOST:="127.0.0.1"}; : ${PORT:=19000}
: ${USER:="default"}; : ${PASS:="default_password"}

if [ -f "/data/db_initialized.flag" ]; then
    s6-svc -d /var/run/s6/services/init-remote-db
    exit 0
fi

bashio::log.info "Waiting for ClickHouse at $HOST:$PORT..."
while ! nc -z "$HOST" "$PORT"; do sleep 2; done

bashio::log.info "Connected! Running initialization scripts..."

# 注入環境變數供 .sh 腳本使用
export CLICKHOUSE_HOST="$HOST"
export CLICKHOUSE_PORT="$PORT"
export CLICKHOUSE_USER="$USER"
export CLICKHOUSE_PASSWORD="$PASS"

for script_file in /app/sql/*; do
    [ -f "$script_file" ] || continue
    chmod +x "$script_file"
    bashio::log.info "Executing: $script_file"
    
    # 使用 bash 執行
    if /bin/bash "$script_file"; then
        bashio::log.info "Success: $script_file"
    else
        bashio::log.error "Failed: $script_file"
    fi
done

touch "/data/db_initialized.flag"
bashio::log.info "Initialization complete."
s6-svc -d /var/run/s6/services/init-remote-db
